# docker/cron/crontab
# m h dom mon dow command
0 */6 * * * php /var/www/lang_app/api/artisan schedule:run
